# my_scripts/monitor_cars/crontab
0 */8 * * * root docker start fetch_cars

