# my_scripts/monitor_cars/crontab
0 * * * * root docker start fetch_cars

